# supervisor/tomcat.conf
[program:tomcat]
user=root
command=/bin/bash /usr/local/content/tomcat/current/bin/catalina.sh run
autostart=true
autorestart=true
redirect_stderr=true
